# contrib/gateway-peers/deploy/agentm-gateway.service
# agentm-gateway.service — the in-memory gateway daemon.
#
# PRIMARY FORM: system unit. Install to /etc/systemd/system/agentm-gateway.service
# and run it as a chosen unprivileged user (see [Service] User=/Group= below).
#
# USER-UNIT VARIANT: if you'd rather run it as `systemctl --user`, install to
# ~/.config/systemd/user/agentm-gateway.service, DROP the User=/Group= lines (a
# user unit already runs as you), and replace the pinned socket with the
# per-user runtime dir, e.g.:
#     RuntimeDirectory=agentm        # → $XDG_RUNTIME_DIR/agentm (user units)
#     ExecStart=… --bind unix://%t/agentm/gw.sock
#     Environment=AGENTM_SOCKET=unix://%t/agentm/gw.sock
# (%t expands to $XDG_RUNTIME_DIR for user units, /run for system units.) Then
# `systemctl --user enable --now agentm-gateway`. Everything else is identical.
#
# Every CAPITALISED /REPLACE/… value and User/Group below MUST be edited before
# enabling. See README.md.

[Unit]
Description=AgentM gateway (single-process chat session host)
Documentation=file:///REPLACE/with/path/to/AgentM/.claude/designs/single-process-gateway.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Run as a dedicated unprivileged account (NOT root). The unix socket's
# peer-cred auth defaults to allowing only this same uid, so the feishu client
# unit MUST run as the same User=/Group= for the connection to be accepted.
User=REPLACE_WITH_SERVICE_USER
Group=REPLACE_WITH_SERVICE_GROUP

# The repo checkout. `uv run` resolves the workspace from here.
WorkingDirectory=/REPLACE/with/path/to/AgentM

# Repo .env carries LARK_*, AGENTM_*, and model credentials. Leading `-` makes
# a missing file non-fatal (the unit still starts). Adjust the path to match
# WorkingDirectory.
EnvironmentFile=-/REPLACE/with/path/to/AgentM/.env

# Pin the socket explicitly. default_socket_url() would pick
# $XDG_RUNTIME_DIR/agentm-gw.sock or /tmp/agentm-gw-<uid>.sock, which is
# ambiguous under systemd (no login session ⇒ no XDG_RUNTIME_DIR), so both the
# --bind flag and AGENTM_SOCKET name the same /run/agentm path. The feishu unit
# connects to this exact path.
RuntimeDirectory=agentm
RuntimeDirectoryMode=0750
Environment=AGENTM_SOCKET=unix:///run/agentm/gw.sock

ExecStart=uv run agentm gateway --bind unix:///run/agentm/gw.sock

# Auto-restart on crash; start on boot (via [Install] below + `enable`).
Restart=always
RestartSec=2

# Let `uv run` find the uv binary even with a sparse systemd PATH. Adjust if
# uv lives elsewhere (e.g. ~/.local/bin for a per-user install — see README).
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
